# ci/build/innovatex-ide@.service
[Unit]
Description=innovatex-ide
After=network.target

[Service]
Type=exec
ExecStart=/usr/bin/innovatex-ide
Restart=always
User=%i

[Install]
WantedBy=default.target
