# ci/build/innovatex-ide-user.service
[Unit]
Description=innovatex-ide
After=network.target

[Service]
Type=exec
ExecStart=/usr/bin/innovatex-ide
Restart=always

[Install]
WantedBy=default.target
